# jprochat.service
[Unit]
Description=jprochat service
After=network.target
StartLimitIntervalSec=0
[Service]
Type=simple
Restart=always
RestartSec=3
User=ubuntu
ExecStart=/home/ubuntu/startJprochat.sh

[Install]
WantedBy=multi-user.target
